# files/latest/pg-backup.timer
[Unit]
Description=Take backups of databases daily

[Timer]
# Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
