# files/latest/pg-backup-cleanup.timer
[Unit]
Description=Cleanup backups of databases daily

[Timer]
# Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
